# contrib/pkg-cacher-report.timer
[Unit]
Description=pkg-cacher report timer

[Timer]
OnCalendar=hourly
Unit=pkg-cacher-report.service

[Install]
WantedBy=timers.target
